# profile.d/haskell.sh
#!/bin/sh
[ -d ~/.cabal/bin ] && PATH=$PATH:~/.cabal/bin

export PATH
