# profile.d/rust.sh
[ -d ~/.cargo/bin ] && PATH=$PATH:~/.cargo/bin

export PATH
