# profile.d/rust.sh
[ -d ~/.multirust/toolchains/nightly/cargo/bin ] && PATH=$PATH:~/.multirust/toolchains/nightly/cargo/bin
[ -d ~/.multirust/toolchains/beta/cargo/bin ] && PATH=$PATH:~/.multirust/toolchains/nightly/cargo/bin
[ -d ~/.multirust/toolchains/stable/cargo/bin ] && PATH=$PATH:~/.multirust/toolchains/nightly/cargo/bin

export PATH
