# profile.d/texlive.sh
#!/bin/sh
[ -d /opt/texlive/2016/bin/x86_64-linux ] && PATH=$PATH:/opt/texlive/2016/bin/x86_64-linux

export PATH
